# home/dot_config/systemd/user/rangerbot.service
[Unit]
Description=Feature-full Telegram bot with triggers capability
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
WorkingDirectory=/home/sonico/Documents/Dev/Bots/NeoKimberly_prod
ExecStart=/home/sonico/.virtualenvs/rangerbot/bin/python /home/sonico/Documents/Dev/Bots/NeoKimberly_prod/kimberly
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
