# root/etc/systemd/system/mnt-share.automount
[Unit]
Description=Automount public share

[Automount]
Where=/mnt/share
TimeoutIdleSec=340

[Install]
WantedBy=multi-user.target
